# debian_files/kvsynopd.cron.d
#Check that this machine owns the ip-alias for kvalobs or if this is a test machine.
#On a test machine the file $HOME/etc/KVALOBS_TEST exist.

25 5 * * * kvalobs /usr/bin/kvsynopdbadmin %
